# init/corosync.service.in
[Unit]
Description=Corosync Cluster Engine
ConditionKernelCommandLine=!nocluster
Requires=network-online.target
After=network-online.target

[Service]
ExecStartPost=/usr/bin/systemctl start pacemaker.service
ExecStop=/usr/bin/systemctl stop pacemaker.service

# If you use watchdog of corosync, uncomment the line below.
#ExecStartPre=/sbin/modprobe softdog soft_margin=60

ExecStart=@INITWRAPPERSDIR@/corosync start
ExecStop=@INITWRAPPERSDIR@/corosync stop
Type=forking

[Install]
WantedBy=multi-user.target
